# infrastructure/ansible/roles/prepare/files/kkkjbot.service
[Unit]
Description=kkkj Discord Bot

[Service]
User=bot
WorkingDirectory=/home/bot/kkkj
ExecStart=/home/bot/kkkj/bot
Restart=always
RestartSec=30